# oe/recipes/powervr-drivers/omap4-sgx-modules-1.7.9.2.1.1/rc.pvr
#!/bin/sh
### BEGIN INIT INFO
# Provides:          pvr-init
# Required-Start:    checkroot
# Required-Stop:
# Default-Start:     S
# Default-Stop:
# Short-Description: PVR SGX driver init.
### END INIT INFO

[ -e /sbin/modprobe ] || exit 1

LOAD_MODULE=/sbin/modprobe

[ "$VERBOSE" != no ] && echo -n "Loading PVR modules: "

kernel=`uname -r | cut -c 1-8`
if [ "$kernel" == "2.6.35.3" ]; then
	exit 0
fi

modprobe omapdrm_pvr
